# deploy/systemd/seo-command-center-worker.service
[Unit]
Description=SEO Command Center background worker
After=network.target

[Service]
# Replace "deploy" with your actual Linux app user.
User=deploy
Group=deploy
WorkingDirectory=/var/www/seo-command-center
Environment=NODE_ENV=production
ExecStart=/bin/bash /var/www/seo-command-center/deploy/run-worker.sh
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
